# templates/command/supervisor.conf.j2
[program:{{ item.key }}]
command={{ item.value['command'] }}
autostart=true
autorestart=true
startretries=20
stdout_logfile=/var/log/{{ item.key }}.log
redirect_stderr=true
