# templates/supervisor_command.conf.j2
[program:{{ name }}]
command={{ command }}
autostart=true
autorestart=true
startretries=20
stdout_logfile=/var/log/{{ name }}.log
redirect_stderr=true
